# namegenerator.service
[Unit]
Description=namegenerator
[Service]
TimeoutStartSec=0
ExecStart=/bin/java -Xms160m -Xmx160m -Xmn70m -jar /root/NameGenerator/build/libs/namegenerator-1.0.jar --server.port=5002
[Install]
WantedBy=multi-user.target
